# scripts/rose/llm-benchrun.service
[Unit]
Description=llm-bench rose-native benchmark run (local exec, resumable)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/llm-bench
Environment=BENCH_LOCAL=1
# The run stops the systemd llama-server router (passwordless sudo) while it holds the GPU,
# and restarts it on exit. Local exec: host scripts run here, not over SSH.
ExecStart=/usr/bin/env node runners/bench-run.mjs --local --resume \
  --models Qwen3.6 \
  --benches triage,reasoning,reasoning_hard,reasoning_expert,toolcalling,summarization,docqa,coding_hard,coding_practical,coding_bugfix,agentic_loop,struct_output,instruction_following \
  --think both --samples 1 --ctx 40960
# Finite, idempotent job: on crash restart and --resume continues; on clean completion it
# exits 0 and stays inactive until the next boot (linger) re-checks and no-ops.
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
